# backend/project-agent-backend.service
[Unit]
Description=Project Agent Backend API
After=network.target postgresql.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/.openclaw/workspace/project-agent/backend
Environment="PATH=/home/ubuntu/.openclaw/workspace/project-agent/backend/venv/bin:/usr/local/bin:/usr/bin:/bin"
# 从 .env 文件加载环境变量
EnvironmentFile=/home/ubuntu/.openclaw/workspace/project-agent/backend/.env
ExecStart=/home/ubuntu/.openclaw/workspace/project-agent/backend/venv/bin/python -m uvicorn app.main:app --host 0.0.0.0 --port 3000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target